# cron/crontab
30 */4 * * * /local/src/docweb/cron/daily/populatedocs
 